# effective_mobile/test-process-monitor.service
[Unit]
Description=Мониторинг процесса test и отправка статуса на API
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitor-test-process.sh
User=root
Group=root
WorkingDirectory=/
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
